# deploy/pi/ar-local-daily.timer
[Unit]
Description=Run AR-local CDR ingest/export daily

[Timer]
# 01:00 Pi local time (Australia/Hobart on ar-local-pi5; see ar_local_ingest_schedule.py).
OnCalendar=*-*-* 01:00:00
Persistent=true
Unit=ar-local-daily.service

[Install]
WantedBy=timers.target
