# deploy/pi/ar-local-daily-watchdog.timer
[Unit]
Description=Check AR-local daily ingest completion

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Persistent=true
Unit=ar-local-daily-watchdog.service

[Install]
WantedBy=timers.target
