# deploy/pi/ar-local-daily.timer
[Unit]
Description=Run AR-local CDR ingest/export daily

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
Unit=ar-local-daily.service

[Install]
WantedBy=timers.target
