# deploy/pi/ar-local-daily-watchdog.service
[Unit]
Description=AR-local daily ingest watchdog catch-up
After=network-online.target
Wants=network-online.target
OnFailure=ar-local-ingest-alert.service

[Service]
Type=oneshot
User={{AR_LOCAL_USER}}
Group={{AR_LOCAL_GROUP}}
WorkingDirectory={{AR_LOCAL_REPO}}
Environment=AR_LOCAL_PORTABLE_ROOT={{AR_LOCAL_PORTABLE_ROOT}}
Environment=AR_LOCAL_DATA_ROOT={{AR_LOCAL_DATA_ROOT}}
# Optional mobile-app payload publishing (inherited by the pi_daily_sync subprocess).
EnvironmentFile=-/etc/ar-local/app-payload.env
ExecStart=/usr/bin/python3 {{AR_LOCAL_REPO}}/pi_daily_watchdog.py
