# deploy/pi/ar-local-deploy-watchdog.service
[Unit]
Description=AR-local Pi deploy verify (sync + /api/latest smoke)
After=network-online.target ar-local-dashboard.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{AR_LOCAL_REPO}}
Environment=AR_PI_VERIFY_LOCAL=1
Environment=AR_PI_BASE_URL=http://127.0.0.1/
ExecStart={{AR_LOCAL_REPO}}/deploy/pi/ar-local-deploy-watchdog.sh

[Install]
WantedBy=multi-user.target
